# deploy/systemd/pi-keyboard.service
[Unit]
Description=Pi BLE HID Keyboard
After=bluetooth.target
Requires=bluetooth.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-ble-keyboard
ExecStart=/usr/bin/python3 -m src.btkeyboard.server
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
